# templates/messup-nft.service.j2
# Managed by Ansible — do not edit manually
[Unit]
Description=messup nft port-forward (oneshot apply)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-{{ nft_install_dir }}/env
ExecStart={{ nft_install_dir }}/apply.sh
ExecStop=/bin/sh -c 'nft list table ip forward >/dev/null 2>&1 && nft delete table ip forward || true'

[Install]
WantedBy=multi-user.target
