# templates/nft-messup.openrc.j2
#!/sbin/openrc-run
# Managed by Ansible — thin wrapper around apply.sh

description="messup nft port-forward (apply.sh)"
apply_script="{{ nft_install_dir }}/apply.sh"
mappings_file="{{ nft_install_dir }}/mappings.txt"

depend() {
	need net
	after firewall
}

start() {
	ebegin "Applying messup nft rules via apply.sh"
	if [ ! -x "${apply_script}" ]; then
		eend 1 "missing: ${apply_script}"
		return 1
	fi
	if [ ! -f "${mappings_file}" ]; then
		eend 1 "missing: ${mappings_file}"
		return 1
	fi
	export IN_IF="{{ _nft_in_if }}"
	export OUT_IF="{{ _nft_out_if }}"
	export CFG="${mappings_file}"
	export LIMIT_RATE="{{ nft_limit_rate }}"
	export MARK="{{ nft_mark }}"
	export TC_ENABLED="{{ 'yes' if (nft_tc_enabled | bool) else 'no' }}"
	"${apply_script}"
	eend $?
}

stop() {
	ebegin "Removing nft table ip forward"
	nft list table ip forward >/dev/null 2>&1 && nft delete table ip forward
	eend 0
	if [ "{{ 'yes' if (nft_tc_enabled | bool) else 'no' }}" = "yes" ]; then
		ebegin "Clearing tc root qdisc"
		[ -n "{{ _nft_out_if }}" ] && tc qdisc del dev "{{ _nft_out_if }}" root 2>/dev/null || true
		if [ "{{ _nft_in_if }}" != "{{ _nft_out_if }}" ]; then
			tc qdisc del dev "{{ _nft_in_if }}" root 2>/dev/null || true
		fi
		eend 0
	fi
}

status() {
	if nft list table ip forward >/dev/null 2>&1; then
		einfo "table ip forward present"
		return 0
	fi
	einfo "table ip forward absent"
	return 3
}
